# luci-app-campus-network-login/root/etc/uci-defaults/luci-campus-network-login
#!/bin/sh

chmod +x /usr/libexec/campus_login

exit 0